# root/etc/services.d/sshd/run
#!/usr/bin/with-contenv sh

# Start sshd
exec /usr/sbin/sshd -f /etc/ssh/sshd_config
